# cups/rootfs/etc/services.d/dbus/run
#!/usr/bin/with-contenv sh
exec /usr/bin/dbus-daemon --system --nofork
